# supervisord.conf
[supervisord]
nodaemon=true

[program:puma]
directory=/app
command=/usr/local/bundle/bin/bundle exec puma -C config/puma.rb

[program:cron]
command=/usr/sbin/crond -f
